# src/php/supervisor_app.conf
[program:php-fpm]
command=php-fpm

[program:sshd]
command=/usr/sbin/sshd -D

[program:cron]
autorestart=false
command=cron -f

[program:postfix]
command=/etc/postfix/postfix-init.sh